# scripts/launchd/com.fireframe.server.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fireframe.server</string>

    <key>ProgramArguments</key>
    <array>
        <!-- Replace with the absolute path to scripts/start.sh on your Mac -->
        <string>/Users/YOUR_USERNAME/path/to/FireFrame/scripts/start.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/fireframe.out</string>

    <key>StandardErrorPath</key>
    <string>/tmp/fireframe.err</string>
</dict>
</plist>
